# scripts/usr/bin/rpi-resize2fs
#!/bin/sh
### BEGIN INIT INFO
# Provides:          rpi-resize2fs
# Required-Start:
# Required-Stop:
# Default-Start: 2 3 4 5 S
# Default-Stop:
# Short-Description: Resize the root filesystem to fill partition
# Description:
### END INIT INFO
. /lib/lsb/init-functions
case "$1" in
  start)
    log_daemon_msg "Starting rpi-resize2fs" &&
    resize2fs /dev/mmcblk0p2 &&
    service rpi-resize2fs.service remove &&
    log_end_msg $?
    ;;
  *)
    echo "Usage: $0 start" >&2
    exit 3
    ;;
esac
